# install/sfn/sfn.service
[Unit]
Description="SafeNetworking Service"

[Service]
ExecStartPre=/bin/sleep 60
ExecStart=/usr/local/bin/sfn.sh

[Install]
WantedBy=default.target
